# ftpd/ftpd.conf
[program:ftpd]
command=/usr/sbin/vsftpd
stdout_logfile=/opt/app/log/%(program_name)s.log
stderr_logfile=/opt/app/log/%(program_name)s.log
autorestart=true
