# ftpd/ftpd.conf
[program:ftpd]
command="/usr/bin/python -m pyftpdlib -w -p 8921 -d /home/tmp"
stdout_logfile=/opt/app/log/%(program_name)s.log
stderr_logfile=/opt/app/log/%(program_name)s.log
autorestart=true
