# orientdb/orientdb.conf
[program:orientdb]
command=/usr/local/bin/server.sh
stdout_logfile=/opt/app/log/%(program_name)s.log
stderr_logfile=/opt/app/log/%(program_name)s.log
autorestart=true
