# levelhud/levelhud.conf
[program:multilevel]
command=/opt/app/index.js
stdout_logfile=/opt/app/log/%(program_name)s.log
stderr_logfile=/opt/app/log/%(program_name)s.log
autorestart=true
