# python-example/app.conf
[program:pandora-local]
command=/opt/app/bin/runserver
stdout_logfile=/opt/app/log/%(program_name)s.log
stderr_logfile=/opt/app/log/%(program_name)s.log
autorestart=true
