# python-example/app.conf
[program:app]
command=/opt/.venv/app/bin/python /opt/app/main.py
stdout_logfile=/opt/app/log/%(program_name)s.log
stderr_logfile=/opt/app/log/%(program_name)s.log
autorestart=true
