# maria/app.conf
[program:mariadb]
command=/usr/bin/mysqld
stdout_logfile=/opt/app/log/%(program_name)s.log
stderr_logfile=/opt/app/log/%(program_name)s.log
autorestart=true
