# maria/maria.conf
[program:mariadb]
command=/opt/maria.sh
stdout_logfile=/opt/app/log/%(program_name)s.out.log
stderr_logfile=/opt/app/log/%(program_name)s.err.log
autorestart=true
