# memsql/memsql.conf
[program:memsql]
command=/usr/local/bin/memsqld
stdout_logfile=/opt/app/log/%(program_name)s.log
stderr_logfile=/opt/app/log/%(program_name)s.log
autorestart=true
