# common/common.conf
[program:sshd]
command=/usr/sbin/sshd -D
stdout_logfile=/opt/app/log/%(program_name)s.log
stderr_logfile=/opt/app/log/%(program_name)s.log
autorestart=true
